# dockerfs/etc/services.d/epson_projector/run
#!/command/with-contenv bash

python /usr/src/app/asyncio_mqtt_based_client.py
